# backup/supervisord_adcuratio.conf
[program:flower]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio flower"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/flower_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/flower_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopsignal=KILL
[program:celery_beat]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery beat -A adcuratio --loglevel=INFO"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/celery_beat_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/celery_beat_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:generate_bit_vector]
command=bash -c "export C_FORCE_ROOT=true && export GMAIL_PASSWORD_DISABLE='adcuratio123' && SENTRY_DSN_DISABLE='https://f38f4ff85a4fbe91a9e9d9dc7483:1cde6eb0f11c481180080d32131abf45@sentry.io/157476' && export SENDGRID_API_KEY_DISABLE='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' && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend && python3 manage.py generate_bit_vector"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/generate_bit_vector_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/generate_bit_vector_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:s3_task_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q s3_task_queue --loglevel=DEBUG --concurrency=1 -n s3_task"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/s3_task_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/s3_task_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:mail_task_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export GMAIL_PASSWORD_DISABLE='adcuratio123' && SENTRY_DSN_DISABLE='https://f7338f4ff85a4f12be91a9e9d9dc7483:1cde6eb0f11c481180080d32131abf45@sentry.io/157476' && export SENDGRID_API_KEY_DISABLE='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' && export OPERATOR_ID=1  && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q mail_task_queue --loglevel=DEBUG --concurrency=1 -n mail"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/mail_task_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/mail_task_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:approve_trade_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 &&  source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q approve_trade_queue --loglevel=DEBUG --concurrency=1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/approve_trade_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/approve_trade_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:populate_edi_info_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 &&  source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q populate_edi_info_queue --loglevel=DEBUG --concurrency=1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/populate_edi_info_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/populate_edi_info_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:pupulate_dummy_and_create_order_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 &&  source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q pupulate_dummy_and_create_order_queue --loglevel=DEBUG --concurrency=1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/pupulate_dummy_and_create_order_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/pupulate_dummy_and_create_order_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:primary_task_queue1]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=SLAVE1 && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q primary_task_queue1 --loglevel=DEBUG --concurrency=1 -n pq1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/primary_task_queue1_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/primary_task_queue1_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:bitvector_task_queue1]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=SLAVE1 && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q bitvector_task_queue1 --loglevel=DEBUG --concurrency=1 -n bq1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/bitvector_task_queue1_id_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/bitvector_task_queue1_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:primary_task_queue2]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=SLAVE2 && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q primary_task_queue2 --loglevel=DEBUG --concurrency=1 -n pq2"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/primary_task_queue2_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/primary_task_queue2_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopsignal=KILL
[program:bitvector_task_queue2]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=SLAVE2 && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q bitvector_task_queue2 --loglevel=DEBUG --concurrency=1 -n bq2"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/bitvector_task_queue2_id_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/bitvector_task_queue2_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopsignal=KILL
[program:miscellaneous_task_queue1]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=SLAVE1 && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q miscellaneous_task_queue1 --loglevel=DEBUG --pool=gevent --concurrency=100 -n mq1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/miscellaneous_task_queue1_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/miscellaneous_task_queue1_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopsignal=KILL
[program:miscellaneous_task_queue2]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=SLAVE2 && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q miscellaneous_task_queue2 --loglevel=DEBUG --pool=gevent --concurrency=100 -n mq2"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/miscellaneous_task_queue2_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/miscellaneous_task_queue2_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopsignal=KILL
[program:celery_beat_worker]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q celery_beat_worker --loglevel=DEBUG --concurrency=1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/celery_beat_worker_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/celery_beat_worker_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:create_icdx_files]
command=bash -c "export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend && python3 manage.py process_pending_campaigns"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/process_pending_campaigns_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/process_pending_campaigns_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:get_or_create_orderline_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 &&  source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q get_or_create_orderline_queue --loglevel=DEBUG --concurrency=1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/get_or_create_orderline_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/get_or_create_orderline_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:pusher_task]
command=python3 -u "/home/yashas/adcuratio/core-backend/pusher_xandr.py"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/pusher_task_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/pusher_task_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
[program:report_task_queue]
command=bash -c "export C_FORCE_ROOT=true && export SERVER_NAME=MASTER && export DEBUG=true && export OPERATOR_ID=1 && source /home/yashas/adcuratio/venv-core-backend/bin/activate && cd /home/yashas/adcuratio/core-backend/ && /home/yashas/adcuratio/venv-core-backend/bin/celery -A adcuratio worker -Q report_task_queue --loglevel=INFO --concurrency=1"
stdout_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/report_task_queue_out.log
stderr_logfile=/home/yashas/adcuratio/adcuratio_v2/logs/report_task_queue_err.log
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
